# debian/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop agent-vcpu.service || true
systemctl disable agent-vcpu.service || true
